# docker/hivy/supervisor-hivy.conf
[supervisord]
nodaemon=true

[program:hivy]
command=hivy --bind 0.0.0.0 --debug

[program:sshd]
command=/usr/sbin/sshd -D

[program:serf]
command=/usr/local/bin/startup-serf
