# docker/node/supervisor-node.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:salt-minion]
command=/usr/bin/startup-salt-minion

[program:serf]
command=/usr/bin/startup-serf
